# ansible/roles/trusttunnel/templates/trusttunnel.service.j2
# trusttunnel.service — systemd unit for the TrustTunnel endpoint daemon
# Generated by Ansible. Do not edit manually — changes will be overwritten.
# Source: roles/trusttunnel/templates/trusttunnel.service.j2

[Unit]
Description=TrustTunnel VPN Endpoint
Documentation=https://github.com/TrustTunnel/TrustTunnel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{ trusttunnel_binary }} {{ trusttunnel_vpn_config }} {{ trusttunnel_hosts_config }}
Restart=on-failure
RestartSec=5s

# Hardening — limit what the process can do if compromised
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ trusttunnel_install_dir }}
# Port 443 requires cap_net_bind_service for non-root
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
